# deploy/arcentra-agent.service
[Unit]
Description=Arcentra CI/CD Platform Agent
Documentation=https://github.com/arcentrix/arcentra
After=network.target
Wants=network.target

[Service]
Type=simple
User=arcentra-agent
Group=arcentra-agent
WorkingDirectory=/opt/arcentra-agent
ExecStart=/opt/arcentra-agent/arcentra-agent -conf /opt/arcentra-agent/conf.d/agent.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arcentra-agent

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/arcentra-agent

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 环境变量
Environment="GOMAXPROCS=0"
Environment="TZ=Asia/Shanghai"

[Install]
WantedBy=multi-user.target
