# deploy/arcentra.service
[Unit]
Description=Arcentra CI/CD Platform Server
Documentation=https://github.com/arcentrix/arcentra
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=arcentra
Group=arcentra
WorkingDirectory=/opt/arcentra
ExecStart=/opt/arcentra/arcentra -conf /opt/arcentra/conf.d/config.toml -plugin /opt/arcentra/plugins
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arcentra

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/arcentra

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 环境变量
Environment="GOMAXPROCS=0"
Environment="TZ=Asia/Shanghai"

[Install]
WantedBy=multi-user.target
